# templates/varnish.service.erb
[Service]
Type=simple

ExecStart=
ExecStart=/usr/sbin/varnishd \
  -F \
  -a :6081 \
  -a localhost:8443,PROXY \
  -T 127.0.0.1:6082 \
  -p feature=+http2 \
  -f /etc/varnish/default.vcl \
  -s file,/var/lib/varnish/api.cache,35G

  ExecReload=/usr/share/varnish/varnishreload -n /var/lib/varnish/<%= node['fqdn'] %>